# deploy/cudackpt@.service
[Unit]
Description=cudackpt control helper (%i)
PartOf=cudackpt.socket
After=cudackpt-run.service cudackpt.socket

[Service]
Type=oneshot
User=root
Group=cudackpt
Environment=CUDACKPT_RUN_DIR=/run/cudackpt
ExecStart=/usr/bin/cudackpt serve
StandardInput=socket
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
